# scripts/com.cmpak.telecom-news-fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cmpak.telecom-news-fetch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/huwenxiao/Downloads/For Claude/pta-telecom-dashboard/scripts/run_news_fetch.sh</string>
    </array>

    <!-- 每天 09:30 PKT 抓取新闻、生成摘要、commit+push 到 GitHub，供 10:10 digest 使用 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/telecom_news_fetch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/telecom_news_fetch_error.log</string>
</dict>
</plist>
