# scripts/com.cmpak.telecom-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cmpak.telecom-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/huwenxiao/Downloads/For Claude/pta-telecom-dashboard/scripts/run_digest.sh</string>
    </array>

    <!-- 每天早上 8:00 触发（PKT 本地时间）-->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Mac 睡眠期间错过的任务，唤醒后立即补跑 -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/telecom_digest.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/telecom_digest_error.log</string>
</dict>
</plist>
